# etc/systemd/system/host_info.service
[Unit]
Description=info for host
After=network-online.target

[Service]
ExecStart=/home/host_info.sh -s

