# etc/systemd/system/host_info.service
[Unit]
Description=info for host
After=network-online.target

[Service]
User=root
Group=root
ExecStart=/usr/bin/host_info.sh -s

